# postgres_installs/etc/profile.d/postgresql.sh
export PATH=${PATH}:/usr/local/postgresql/bin/:/usr/local/pgbackrest/bin/
